# sweethome3d/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting PHP-FPM..."
exec php-fpm82 --nodaemonize --fpm-config /etc/php82/php-fpm.conf
