# sweethome3d/rootfs/etc/services.d/nginx/finish
#!/bin/sh
exec s6-svscanctl -t /run/s6/services
